# attic/scripts/example.cron
#example of crontab, copied from g2 nearline
#warn if disk space gets too low
*/30 * * * * /home/daq/database/diskSpaceWarning.sh 94 > /home/daq/database/disk.log 2>&1

# automatic running of the nearline analyzer(s)
*/5 * * * * /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/run_nearline_automatic.sh > /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/process.log 2>&1

# hadd the outputs
0,20,40 * * * * source /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/do_hadd.sh > /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/hadd_cron.log 2>&1

*/30 * * * * source /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/midas_backup.sh > /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/midas_backup.log 2>&1

0 0,12 * * * source /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/psql_db_backup.sh > /home/pioneer/pioneer/pioneer-teststand/packages/unpacking/scripts/psql_db_backup.log 2>&1
